# puppet/modules/site-linode/files/etc/cron.d/backup-gdbm

0 8 * * * git /usr/local/bin/backup-gdbm.sh
